# packages/letux-gadget/core/etc/udev/rules.d/gadget.rules
# Rule for setting up the gadget (configfs) drivers after boot

ACTION=="*", SUBSYSTEM=="gadget", RUN+="/etc/udev/rules.d/gadget.sh"
