# packages/letux-gadget/ip-over-usb/etc/init.d/letux-ip-over-usb-gadget@usb0
#!/bin/sh
# kFreeBSD do not accept scripts as interpreters, using #!/bin/sh and sourcing.
if [ true != "$INIT_D_SCRIPT_SOURCED" ] ; then
    set "$0" "$@"; INIT_D_SCRIPT_SOURCED=true . /lib/init/init-d-script
fi
### BEGIN INIT INFO
# Provides:          letux-ip-over-usb-gadget
# Required-Start:
# Required-Stop:
# Default-Start:     1 2 3 4 5
# Default-Stop:      0 6
# Short-Description: Control IP over USB (ncm)
# Description:       This script adds usb0 and a dhcp daemon to configfs
#                    
### END INIT INFO

# Author: H. Nikolaus Schaller <hns@goldelico.com>

SERVICE=letux-ip-over-usb-gadget

do_status () {
	if [ -f /sys/class/net/usb0 ] ; then
		echo "$SERVICE is running."
		return 0
	else
		return 4
	fi
}

case "$1" in
  start)
        log_daemon_msg "Starting $SERVICE"
	/etc/gadget/core/letux-gadget.sh add ncm
	/sbin/ifconfig usb0 netmask 255.255.255.192 192.168.0.202 up
	/sbin/route add default gw 192.168.0.200 metric 203 | : may fail
	/etc/network/letux-setup-resolv-conf.sh
# FIXME: start dhcpd service and just copy configuration file in the right place...
	/usr/sbin/dhcpd -q -cf /etc/gadget/scripts/letux-ip-over-usb-gadget-dhcpd.conf
        exit 0
        ;;
  restart|reload|force-reload)
        echo "Error: argument '$1' not supported" >&2
        exit 3
        ;;
  stop)
        log_daemon_msg "Stopping $SERVICE"
# FIXME: stop dhcpd
	ifconfig usb0 down || : ignore error
	/etc/gadget/core/letux-gadget.sh remove ncm
        log_progress_msg "disabled."
        log_end_msg 0
        exit 0
        ;;
  status)
	do_status
	exit $?
	;;
  *)
        echo "Usage: $0 start|stop|status" >&2
        exit 3
        ;;
esac

:
